# tools/launchd/com.user.claude-watchlist-refresh.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template — install_watchlist_refresh.sh substitutes __REPO_PATH__ and
  __PYTHON_BIN__, writes the rendered plist to ~/Library/LaunchAgents/,
  then loads it. Edit this template to change schedule
  (StartCalendarInterval) or paths; re-run the install script to apply.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.claude-watchlist-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>cd "__REPO_PATH__" &amp;&amp; "__PYTHON_BIN__" "__REPO_PATH__/tools/refresh_watchlist.py"</string>
  </array>
  <!-- Why script path not `-m tools.refresh_watchlist`: under launchd's
       minimal shell environment, `python3 -m` resolution from cwd is
       unreliable when the repo path contains special characters
       (e.g. OneDrive's CJK path components). Invoking the script
       directly sidesteps the issue — the script's own
       sys.path.insert(0, ROOT) handles the imports. -->

  <!-- Ensure PATH includes /usr/bin and any user-bin locations the
       script may transitively shell out to (currently none, but
       defensive against future additions). -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- Two weekday runs:
       07:00 — pre-market warm-up: refreshes the prior completed trading
       day + overnight US closes so the dashboard loads instantly at open.
       17:00 — post-close: this is the trading-day rollover hour (see
       TRADING_DAY_ROLLOVER_HOUR in app.py), so by now today's daily
       OHLCV / T86 / 融資借券 are published and the cache tag has advanced
       to today — this run pulls today's complete data.
       launchd auto-catches up if the Mac was asleep (unlike cron). -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <!-- Persistent log so cron diagnosis works across reboots. -->
  <key>StandardOutPath</key>
  <string>__REPO_PATH__/tools/launchd/watchlist-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_PATH__/tools/launchd/watchlist-refresh.log</string>

  <!-- Don't restart if it exits (one-shot per schedule trigger). -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
