# tools/launchd/com.user.claude-broker-branch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template — install_broker_branch.sh substitutes __REPO_PATH__ and
  __PYTHON_BIN__, writes the rendered plist to ~/Library/LaunchAgents/,
  then loads it. Edit this template to change schedule / paths; re-run the
  install script to apply.

  Headless run of broker_branch_pull.py --no-crawl: the user prepares
  step 2 (the BSR crawl) earlier in the day, so by 17:00 the CSVs are
  already in twse_web/output/ and this job only INGESTS them (no crawl, no
  interactive pause, no terminal). The crawl — which solves a CAPTCHA — is
  never run or scheduled by this job; it's the user's manual prep.
  Runs at 17:00, AFTER the 16:30 watchlist-refresh, so today's per-stock
  `lots` are already cached and the volume-based date auto-detect can
  verify each file's day instead of falling back.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.claude-broker-branch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>cd "__REPO_PATH__" &amp;&amp; "__PYTHON_BIN__" "__REPO_PATH__/tools/broker_branch_pull.py" --no-crawl</string>
  </array>
  <!-- Script path not `-m`: under launchd's minimal shell, `-m` resolution
       from a CJK/CloudStorage cwd is unreliable; the script's own
       sys.path handling (via stock_web import) works regardless. -->

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- Weekdays 17:00 — after the 16:30 watchlist-refresh; the user
       crawls the day's 分點 beforehand. launchd auto-catches up if the
       Mac slept. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>__REPO_PATH__/tools/launchd/broker-branch.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_PATH__/tools/launchd/broker-branch.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
